# server/src/packaging/freebsd/onCourseServer.sh.in
#!/bin/sh

# PROVIDE: onCourseServer
# REQUIRE: NETWORKING SERVERS mysql
# KEYWORD: shutdown

#
# Add the following lines to /etc/rc.conf to enable ish onCourse Server:
#
# onCourseServer_enable="YES"
#
# If you find that your onCourseServer is having memory problems:
# set onCourseServer_mem to run with additional memory in /etc/rc.conf
#
#
# For additional options for onCourseServer
# onCourseServer_flags in /etc/rc.conf

. /etc/rc.subr

name="onCourseServer"
rcvar=`set_rcvar`

load_rc_config "${name}"

: ${onCourseServer_enable="NO"}
: ${onCourseServer_flags=""}
: ${onCourseServer_mem="750M"}
: ${onCourseServer_home="%%PREFIX%%/onCourse"}
: ${onCourseServer_jmx_port=""}
: ${onCourseServer_java_home="%%PREFIX%%/openjdk8"}
: ${onCourseServer_log_dir="%%PREFIX%%/onCourse/logs"}

command=/usr/sbin/daemon
pidfile=/var/run/onCourse/onCourseServer.pid
java_cmd="${onCourseServer_java_home}/bin/java"
procname="${java_cmd}"
required_files="${java_cmd}"
onCourseServer_log_file=${onCourseServer_log_dir}/stdout.txt

# We could reduce permissions
onCourse_user=root
onCourse_group=wheel

start_precmd="onCourseServer_prestart"

if [ -n "$onCourseServer_jmx_port" ]; then
   onCourseServer_flags="${onCourseServer_flags} \
     -Dcom.sun.management.jmxremote \
     -Dcom.sun.management.jmxremote.port=${onCourseServer_jmx_port} \
     -Dcom.sun.management.jmxremote.ssl=false \
     -Dcom.sun.management.jmxremote.authenticate=false"
fi

onCourseServer_flags="-p ${pidfile} ${java_cmd} \
   -Xmx${onCourseServer_mem} \
   -server \
   -Djava.awt.headless=true \
   ${onCourseServer_flags} \
   -jar ${onCourseServer_home}/server-%%PORTVERSION%%.jar \
   HEADLESS=true \
   > ${onCourseServer_log_file} 2>&1"

onCourseServer_prestart() {
        if [ ! -f "${onCourseServer_log_file}" ]; then
                install -d -o "${onCourse_user}" -g "${onCourse_group}" -m 750 "${onCourseServer_log_dir}"
                touch "${onCourseServer_log_file}"
                chown "${onCourse_user}:${onCourse_group}" "${onCourseServer_log_file}"
                chmod 640 "${onCourseServer_log_file}"
        fi
        if [ ! -d "/var/run/onCourse" ]; then
                install -d -o "${onCourse_user}" -g "${onCourse_group}" -m 750 "/var/run/onCourse"
        fi
}

cd ${onCourseServer_home}

run_rc_command "$1"